# services/networkmanager-dispatcher/run
#!/bin/bash
set -euo pipefail

node /etc/djs/daemontools.js/svwait.js -u networkmanager

exec /usr/lib/nm-dispatcher
