# services/networkmanager/run
#!/bin/bash
set -euo pipefail

node /etc/djs/daemontools.js/svwait.js -u dbus
node /etc/djs/daemontools.js/svwait.js -u udevadm-oneshot

exec /usr/bin/NetworkManager --no-daemon
