# services/acpid/run
#!/bin/bash
set -euo pipefail

node /etc/djs/daemontools.js/svwait.js -u udevadm-oneshot

exec acpid --foreground --netlink
